# deploy/botpanel.service
[Unit]
Description=BotPanel — bot Discord + API + site web
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=botpanel
Group=botpanel
WorkingDirectory=/opt/botpanel
EnvironmentFile=/opt/botpanel/.env
ExecStart=/opt/botpanel/.venv/bin/python -m app.main
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Securisation (optionnelle)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/botpanel/data
PrivateTmp=true

[Install]
WantedBy=multi-user.target
